# unified-ai-agent/deployment/systemd/unified-agent.service
[Unit]
Description=Unified AI Agent (REPL / service)
After=network.target

[Service]
WorkingDirectory=/opt/unified-ai-agent
Environment=PYTHONPATH=/opt/unified-ai-agent/src
Environment=UNIFIED_PROVIDER=auto
ExecStart=/usr/bin/python3 -m agent.cli --repl
Restart=on-failure
User=unified
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
